# autostart/volumeicon.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=volumeicon
Exec=volumeicon
StartupNotify=false
Terminal=false
Hidden=false
